# panel/systemd/realm-panel.service
[Unit]
Description=Realm Panel Web UI
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=-/etc/realm-panel/panel.env
WorkingDirectory=/opt/realm-panel/panel
ExecStart=/bin/bash -lc 'set -e; if [[ -x /opt/realm-panel/panel/start.sh ]]; then exec /opt/realm-panel/panel/start.sh; fi; if [[ -x /opt/realm-panel/start.sh ]]; then exec /opt/realm-panel/start.sh; fi; echo "panel start script missing under /opt/realm-panel" >&2; exit 1'
Restart=always
RestartSec=2
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
